# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:pentaho]
command=/opt/pentaho/biserver-ce/start-pentaho.sh
#command=service pentaho start



