# supervisord.conf
[supervisord]
nodaemon=true

[program:pentaho]
command=/opt/pentaho/biserver-ce/start-pentaho.sh
#command=service pentaho start



